# files/rc.local
#!/bin/sh 

# turn off tv output
/usr/bin/tvservice --off

